# agent/rc.d/orbit_agent
#!/bin/sh
#
# PROVIDE: orbit_agent
# REQUIRE: NETWORKING
# KEYWORD: shutdown
#
# Add the following line to /etc/rc.conf to enable the agent:
#   orbit_agent_enable="YES"

. /etc/rc.subr

name="orbit_agent"
rcvar="${name}_enable"
desc="OPNsense Dashboard Agent"

: ${orbit_agent_enable:="NO"}
: ${orbit_agent_config:="/usr/local/etc/orbit-agent.conf"}
: ${orbit_agent_pidfile:="/var/run/${name}.pid"}
: ${orbit_agent_logfile:="/var/log/${name}.log"}

pidfile="${orbit_agent_pidfile}"

# Resolve a Python interpreter without hardcoding the version. OPNsense ships
# /usr/local/bin/python3; pfSense may install only a versioned binary (no python3
# symlink) and the version varies by release (3.8 on older boxes, 3.11+ on newer).
# Prefer an unversioned python3, else pick the NEWEST python3.N found.
: ${orbit_agent_python:=""}
if [ -z "${orbit_agent_python}" ]; then
    for _d in /usr/local/bin /usr/bin; do
        if [ -x "${_d}/python3" ]; then orbit_agent_python="${_d}/python3"; break; fi
        _cand=$(ls "${_d}"/python3.* 2>/dev/null | grep -E '/python3\.[0-9]+$' | sort -t. -k2 -rn | head -1)
        if [ -n "${_cand}" ] && [ -x "${_cand}" ]; then orbit_agent_python="${_cand}"; break; fi
    done
fi

# Run the supervisor (run-agent.sh), not python directly: it provides the
# respawn loop + self-update rollback (DR-5). The interpreter is passed as $1.
#
# -f is critical on pfSense: at boot the service is launched by the PHP boot
# framework via afterbootupshellcmd, i.e. mwexec()/exec() which reads the child's
# stdout to EOF. Without -f the long-lived daemon(8) SUPERVISOR keeps its inherited
# stdout/stderr (the exec capture pipe), so EOF never arrives and rc.bootup hangs
# forever at "Running afterbootupshellcmd ... onestart". -f sends the supervisor's
# own fds to /dev/null; combined with -o the agent's output still lands in the log.
command="/usr/sbin/daemon"
command_args="-f -p ${pidfile} -o ${orbit_agent_logfile} -t ${name} /bin/sh /usr/local/orbit-agent/run-agent.sh ${orbit_agent_python}"

start_precmd="${name}_prestart"
# The pidfile holds run-agent.sh's pid, not /usr/sbin/daemon's, so rc.subr's
# default (procname-matching) stop/status can't find it. Drive them off the
# pidfile directly instead.
stop_cmd="${name}_stop"
status_cmd="${name}_status"

orbit_agent_prestart() {
    if [ ! -f "${orbit_agent_config}" ]; then
        echo "Config file ${orbit_agent_config} not found."
        echo "Copy the example: cp /usr/local/orbit-agent/orbit-agent.conf.example ${orbit_agent_config}"
        return 1
    fi
    # Ensure log file exists
    touch "${orbit_agent_logfile}"
}

orbit_agent_stop() {
    if [ -f "${pidfile}" ]; then
        _pid=$(cat "${pidfile}")
        echo "Stopping ${name} (pid ${_pid})."
        # run-agent.sh traps TERM → kills its python child → exits → daemon(8) exits.
        kill -TERM "${_pid}" 2>/dev/null
        wait_for_pids "${_pid}"
    else
        echo "${name} not running? (no pidfile)"
    fi
    rm -f "${pidfile}"
}

orbit_agent_status() {
    if [ -f "${pidfile}" ] && kill -0 "$(cat "${pidfile}")" 2>/dev/null; then
        echo "${name} is running as pid $(cat "${pidfile}")."
    else
        echo "${name} is not running."
    fi
}

load_rc_config $name
run_rc_command "$1"
